# linux/etc/systemd/system/myscript.service

[Unit]
Description=My Startup Script
After=default.target

[Service]
ExecStart=/usr/local/bin/my_script.sh

[Install]
WantedBy=default.target
