# linux/etc/xdg/autostart/script.desktop
[Desktop Entry] Name=Script Type=Application Exec=/usr/local/bin/my_script.sh
