# linux/etc/xdg/autostart/type2color-startup.desktop
[Desktop Entry] Name=Type2Color Startup Type=Application Exec=/usr/local/bin/type2color.sh
